# crontab
0 1 * * * root pip install --upgrade youtube-dl  > /proc/1/fd/1 2>/proc/1/fd/2
*/5 * * * * root bash /app/bash_cron_delete.sh  > /proc/1/fd/1 2>/proc/1/fd/2

